# modules/os/fs/etc/supervisord.d/crond.conf
[program:crond]
command=/usr/sbin/crond -n
priority=3
autorestart=true