# modules/os/fs/etc/supervisord.d/rsyslogd.conf
[program:rsyslogd]
command=/sbin/rsyslogd -n
autorestart=true